# cron/crontab
@daily bash /opt/app-root/cron/backupOdooDB
